# systemd/cat-agents-stability-review.service
[Unit]
Description=Cat Agents Stability One-Time Review Report
After=cat-agents-stabilityd.service openclaw-gateway.service

[Service]
Type=oneshot
User=flashcat
Group=flashcat
WorkingDirectory=/home/flashcat
ExecStart=/home/flashcat/cat-agents-stabilityd/bin/cat-agents-stability-review.sh

